# airootfs/etc/systemd/system/getty@tty7.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin vet --noclear %I 38400 linux
